# docker/service/ivit-i-cli.service
[Unit]
Description=Stop iVIT-I Web Service

[Service]
Type=simple
ExecStart=
ExecStop=

[Install]
WantedBy=default.target